# config/com.theorb.watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.theorb.watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/andrewrooney/Projects/the-orb/.venv/bin/python</string>
        <string>/Users/andrewrooney/Projects/the-orb/src/orb.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/andrewrooney/Projects/the-orb/src</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/andrewrooney/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/andrewrooney</string>
    </dict>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/andrewrooney/Projects/the-orb/logs/launchd.out</string>

    <key>StandardErrorPath</key>
    <string>/Users/andrewrooney/Projects/the-orb/logs/launchd.err</string>

</dict>
</plist>
